# adguard/rootfs/etc/services.d/nginx/finish
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: AdGuard Home
# Take down the S6 supervision tree when Nginx fails
# ==============================================================================
if [[ "${1}" -ne 0 ]] && [[ "${1}" -ne 256 ]]; then
  bashio::log.warning "NGinx crashed, halting add-on"
  exec /run/s6/basedir/bin/halt
fi

bashio::log.info "NGinx stoped, restarting..."
